# crontab
SHELL=/bin/sh 
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

*/5 * * * * root sudo -n -u spray python /app/sprays.py > /proc/1/fd/1 2> /proc/1/fd/2
